# services/164/g1d-offset-detector.service
[Unit]
Description=G1D Lift Offset Detector Service
After=network.target ros2.service
Wants=ros2.service

[Service]
Type=simple
User=__G1D_USER__
Group=__G1D_USER__
WorkingDirectory=__G1D_HOME__
ExecStart=/usr/bin/python3 __G1D_HOME__/bin/g1d_offset_detector.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
