# xiangmu_ceshi_file/services/4090/order-dashboard.service
[Unit]
Description=G1D Order Dashboard (MySQL)
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=__G1D_USER__
Group=__G1D_USER__
WorkingDirectory=__G1D_HOME__
ExecStart=/usr/bin/python3 __G1D_HOME__/bin/order_dashboard.py
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitIntervalSec=60
Environment=MYSQL_HOST=127.0.0.1
Environment=MYSQL_PORT=3306
Environment=MYSQL_USER=bwton
Environment=MYSQL_PASSWORD=bwton@888
Environment=MYSQL_DB=digitaltwins

[Install]
WantedBy=multi-user.target
