# services/5/g1d-control-api.service
[Unit]
Description=G1D Control API Service
After=network.target multi-user.target
Wants=g1d-monitor.service
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User=__G1D_USER__
WorkingDirectory=__G1D_HOME__
# 等30秒让 ROS2 DDS 完成远程节点发现后再启动
ExecStartPre=/bin/sleep 30
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && exec /usr/bin/python3 __G1D_HOME__/bin/g1d_control_api.py'
Restart=on-failure
RestartSec=10
# SIGTERM 后最多等 15 秒，超时则 SIGKILL
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
